# lot/deploy/redis.service
[Unit]
Description=Redis In-Memory Data Store
After=network.target
